# supervisord.conf
[supervisord]
nodaemon=true

[program:app]
command=bun --cwd /app/packages/server /app/packages/server/dist/index.js
directory=/app/packages/server
autostart=true
autorestart=true
